# deploy/com.palmsh.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.palmsh.gateway</string>

	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/node</string>
		<string>__PROJECT_DIR__/packages/gateway/dist/index.js</string>
	</array>

	<key>WorkingDirectory</key>
	<string>__PROJECT_DIR__/packages/gateway</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>__HOME__</string>
		<key>PATH</key>
		<string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>

	<key>UserName</key>
	<string>__USER__</string>

	<key>KeepAlive</key>
	<true/>

	<key>RunAtLoad</key>
	<true/>

	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/palmsh-gateway.log</string>

	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/palmsh-gateway.log</string>

	<key>ThrottleInterval</key>
	<integer>5</integer>
</dict>
</plist>
