# packages/gateway/ccsuite-gateway.service
[Unit]
Description=ccsuite Terminal Gateway
After=network.target
Documentation=https://github.com/your-repo/ccsuite

[Service]
Type=simple
User=raspi5
Group=raspi5
WorkingDirectory=/home/raspi5/projects/ccsuite/packages/gateway
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/home/raspi5/projects/ccsuite/packages/gateway/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ccsuite-gateway

[Install]
WantedBy=multi-user.target
