# packages/gateway/palmsh-gateway.service
[Unit]
Description=zenterm Terminal Gateway
After=network.target
Documentation=https://github.com/your-repo/zenterm

[Service]
Type=simple
User=raspi5
Group=raspi5
WorkingDirectory=/home/raspi5/projects/zenterm/packages/gateway
ExecStart=/usr/bin/node dist/index.js
EnvironmentFile=/home/raspi5/projects/zenterm/packages/gateway/.env
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zenterm-gateway

[Install]
WantedBy=multi-user.target
